# openclaw-voice.service
# OpenClaw Voice - Systemd Service (Linux)
# For other platforms, see SKILL.md for alternatives

[Unit]
Description=OpenClaw Voice Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=%h/openclaw-voice
EnvironmentFile=%h/openclaw-voice/.env
ExecStart=/usr/bin/node %h/openclaw-voice/src/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
